# bundle-templates/systemd/steampipe.service
[Unit]
Description=Steampipe SQL Query Service
Documentation=https://steampipe.io/docs
After=network.target

[Service]
Type=simple
User=steampipe
Group=steampipe

WorkingDirectory=/opt/steampipe
EnvironmentFile=/opt/steampipe/steampipe.env

ExecStart=/opt/steampipe/bin/server
ExecStop=/opt/steampipe/bin/stop

Restart=on-failure
RestartSec=10

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/steampipe

[Install]
WantedBy=multi-user.target
