# config/systemd/steampipe.service
[Unit]
Description=Steampipe SQL Service
Documentation=https://steampipe.io/docs
After=network.target local-fs.target
RequiresMountsFor=/opt/wsl-secrets

[Service]
Type=simple
User=steampipe
Group=steampipe

# Create /run/steampipe as tmpfs
RuntimeDirectory=steampipe
RuntimeDirectoryMode=0755

# Environment
Environment=STEAMPIPE_INSTALL_DIR=/run/steampipe
Environment=STEAMPIPE_MOD_LOCATION=/run/steampipe
Environment=HOME=/opt/steampipe
WorkingDirectory=/opt/steampipe

# Extract and start in one script (no gap for scanner)
ExecStart=/opt/init/steampipe-start.sh

ExecStop=/run/steampipe/steampipe/steampipe service stop

Restart=on-failure
RestartSec=10
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
